# crontabs/root
# min	hour	day	month	weekday	command
*/30	*	*	*	*	run-parts /etc/periodic/30min
0	*	*	*	*	run-parts /etc/periodic/hourly
